# net-misc/networkmanager/files/init.d.NetworkManager-dhcpcd
#!/sbin/openrc-run
# Copyright (c) 2008 Saleem Abdulrasool <compnerd@compnerd.org>
# Copyright 2013-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="NetworkManager daemon. The service is marked as started even \
if a network connection is not established."

depend() {
	need dbus
}

start() {
	ebegin "Starting NetworkManager"
	start-stop-daemon --start --quiet --pidfile /run/NetworkManager/NetworkManager.pid \
		--exec /usr/sbin/NetworkManager -- --pid-file /run/NetworkManager/NetworkManager.pid
	eend $?
}

stop() {
	ebegin "Stopping NetworkManager"
	start-stop-daemon --stop --quiet --pidfile /run/NetworkManager/NetworkManager.pid
	# If dhcpcd is used it must be killed or it will stick around.
	killall -9 dhcpcd || true
	# Starting from 1.2.4 the devices used must be flushed for the service to be restarted.
	# It might actually be restarting but nm-applet itself is glitchy.
	ip addr flush dev eth0 || true
	eend $?
}

# vim: set ft=gentoo-init-d ts=4 :
